# external/packages/bsp/rk356x/etc/udev/rules.d/83-typec.rules
ACTION=="remove", SUBSYSTEM=="typec", RUN+="/usr/local/bin/typec_usb2.sh"
