# external/packages/bsp/rk356x/etc/udev/rules.d/90-load-uwe5622-modules.rules
ACTION=="add", SUBSYSTEM=="sdio", SUBSYSTEMS=="sdio", ENV{SDIO_ID}=="0000:0000", RUN+="/sbin/modprobe -a sprdwl_ng sprdbt_tty"
