# external/packages/bsp/rk3588/etc/udev/rules.d/84-audio-orangepi.rules
ACTION=="change", SUBSYSTEM=="extcon", ATTR{name}=="es8388-sound", RUN+="/usr/bin/orangepi_audio.sh"
