# external/packages/blobs/desktop/skel_new/.config/autostart/org.orangepi.opiwiz.desktop
[Desktop Entry]
Name=opiwiz
Exec=sudo opiwiz
Icon=org.orangepi.opiwiz
Terminal=false
Type=Application
Categories=GTK;
StartupNotify=true
