# external/packages/bsp/rk356x/usr/lib/udev/rules.d/61-usbdevice.rules
SUBSYSTEM=="android_usb",ACTION=="change",RUN+="/usr/bin/usbdevice update"
