# external/packages/bsp/rk3588/etc/profile.d/x11.sh
export DISPLAY=${DISPLAY:-:0}
export TERM=${TERM:-linux}
