# external/packages/bsp/common/etc/cron.d/orangepi-updates
@reboot root /usr/lib/orangepi/orangepi-apt-updates
@daily root /usr/lib/orangepi/orangepi-apt-updates
