# external/packages/bsp/common/lib/systemd/system/orangepi-disable-autologin.service
[Unit]
Description=Disable automated desktop login

[Service]
Type=simple
ExecStart=rm -f /etc/lightdm/lightdm.conf.d/22-orangepi-autologin.conf
TimeoutStopSec= 10
ExecStop=systemctl disable orangepi-disable-autologin.timer
