# external/packages/bsp/common/lib/systemd/system/orangepi-resize-filesystem.service
# Orange Pi resize filesystem service
# Resizes partition and filesystem on first/second boot
# This service may block the boot process for up to 3 minutes

[Unit]
Description=Orange Pi filesystem resize
Before=basic.target
After=sysinit.target local-fs.target
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/orangepi/orangepi-resize-filesystem start
TimeoutStartSec=6min

[Install]
WantedBy=basic.target
