# external/packages/bsp/common/lib/systemd/system/orangepi-hardware-monitor.service
# Orange Pi hardware monitor service
# Detects and sets various features on boot
# This service may block the boot process for up to 2 minutes

[Unit]
Description=Orange Pi hardware monitoring
Before=basic.target
After=sysinit.target local-fs.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/lib/orangepi/orangepi-hardware-monitor start
ExecStop=/usr/lib/orangepi/orangepi-hardware-monitor stop
RemainAfterExit=yes
TimeoutStartSec=2min

[Install]
WantedBy=basic.target
