# external/packages/bsp/rk3588/usr/lib/systemd/system/usbdevice.service
[Unit]
Description=Manage USB device functions
DefaultDependencies=no
After=local-fs.target

[Service]
Type=forking
ExecStart=/usr/bin/usbdevice start
ExecStop=/usr/bin/usbdevice stop

[Install]
WantedBy=sysinit.target
